# supervisord.conf
[supervisord]
nodaemon=true

[program:ossec]
command=/var/ossec/bin/ossec-control start